# templates/certbot-renewal.service.j2
# Certbot automatic renewal service
# Managed by Ansible — DO NOT EDIT MANUALLY
[Unit]
Description=Certbot automatic certificate renewal
Documentation=https://certbot.eff.org/docs/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew --quiet \
  --config-dir {{ certbot_config_dir }} \
  --work-dir {{ certbot_work_dir }} \
  --logs-dir {{ certbot_log_dir }}{% if certbot_renewal_pre_hook | length > 0 %} \
  --pre-hook "{{ certbot_renewal_pre_hook }}"{% endif %}{% if certbot_renewal_post_hook | length > 0 %} \
  --post-hook "{{ certbot_renewal_post_hook }}"{% endif %}{% if certbot_renewal_deploy_hook | length > 0 %} \
  --deploy-hook "{{ certbot_renewal_deploy_hook }}"{% endif %}

PrivateTmp=true
ProtectSystem=full
